# deployment/systemd/nfc-websocket.service
[Unit]
Description=NFC Collection WebSocket Server
After=network.target

[Service]
Type=simple
User=loganrhyne
Group=loganrhyne
SupplementaryGroups=gpio
WorkingDirectory=/home/loganrhyne/nfc-collection/python-services
Environment="PATH=/home/loganrhyne/nfc-collection/python-services/venv/bin:/home/loganrhyne/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="VIRTUAL_ENV=/home/loganrhyne/nfc-collection/python-services/venv"
ExecStart=/home/loganrhyne/nfc-collection/python-services/venv/bin/python server.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target