# deployment/systemd/nfc-dashboard.service
[Unit]
Description=NFC Collection Dashboard Web Server
After=network.target nfc-websocket.service
Wants=nfc-websocket.service

[Service]
Type=simple
User=loganrhyne
Group=loganrhyne
WorkingDirectory=/home/loganrhyne/nfc-collection/dashboard-ui
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 -m http.server 80 --directory build
Restart=always
RestartSec=10

# Run on port 80 (requires capability)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target