# deployment/systemd/nfc-server.service
[Unit]
Description=NFC Collection Server
After=network.target

[Service]
Type=simple
User=loganrhyne
Group=gpio
SupplementaryGroups=spi,i2c

WorkingDirectory=/home/loganrhyne/nfc-collection/python-services
Environment="PATH=/home/loganrhyne/nfc-collection/python-services/venv/bin:/usr/bin:/bin"
Environment="PORT=8000"

ExecStart=/home/loganrhyne/nfc-collection/python-services/venv/bin/python server_clean.py

Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target