# deployment/systemd/nfc-hardware.service
[Unit]
Description=NFC Hardware Service
After=multi-user.target

[Service]
Type=simple
User=loganrhyne
Group=gpio
SupplementaryGroups=spi,i2c

WorkingDirectory=/home/loganrhyne/nfc-collection/python-services
Environment="PATH=/home/loganrhyne/nfc-collection/python-services/venv/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/loganrhyne/nfc-collection/python-services"

ExecStart=/home/loganrhyne/nfc-collection/python-services/venv/bin/python hardware_service.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Hardware access permissions
PrivateDevices=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target